# packages/qdtech-touch-x11/autostart/orangepi-touch-map.desktop
[Desktop Entry]
Type=Application
Name=Map QDtech Touchscreen
Exec=/usr/local/bin/orangepi-touch-map
X-GNOME-Autostart-enabled=false
Hidden=true
