# packages/qdtech-touch-x11/autostart/qdtech-touch-x11.desktop
[Desktop Entry]
Type=Application
Name=QDtech Touch X11 Bridge
Exec=/usr/local/bin/qdtech-touch-x11
X-GNOME-Autostart-enabled=true
Hidden=false
